# etc/init.d/dnsmasq_route
#!/bin/sh /etc/rc.common

START=99

DNSMASQ_ROUTE_CONF=/tmp/dnsmasq.d/dnsmasq_nftset.conf
DNSMASQ_ROUTE_SET=gfw_list
DNSMASQ_ROUTE_NFT_TABLE=dnsmasq_route
DNSMASQ_ROUTE_NFT_CHAIN=prerouting
DNSMASQ_ROUTE_RULE_PREF=29000
# nftset 条目固定保留 24 小时，重复解析不会刷新剩余过期时间。
DNSMASQ_ROUTE_NFT_TIMEOUT_DEFAULT=24h

netmask_to_prefix()
{
  local IFS=.
  local octet prefix=0

  for octet in $1
  do
    case "$octet" in
      255) prefix=$((prefix + 8)) ;;
      254) prefix=$((prefix + 7)) ;;
      252) prefix=$((prefix + 6)) ;;
      248) prefix=$((prefix + 5)) ;;
      240) prefix=$((prefix + 4)) ;;
      224) prefix=$((prefix + 3)) ;;
      192) prefix=$((prefix + 2)) ;;
      128) prefix=$((prefix + 1)) ;;
      0) ;;
      *) return 1 ;;
    esac
  done

  echo "$prefix"
}

get_remote_dev()
{
  ip route get "$1" 2>/dev/null | awk '
    {
      for (i = 1; i <= NF; i++) {
        if ($i == "dev") {
          print $(i + 1)
          exit
        }
      }
    }
  '
}

setup_policy_route()
{
  local table="$1"
  local remote_ip="$2"
  local remote_dev

  remote_dev=$(get_remote_dev "$remote_ip")

  ip rule del pref "$DNSMASQ_ROUTE_RULE_PREF" fwmark "$table" table "$table" 2>/dev/null
  ip rule add pref "$DNSMASQ_ROUTE_RULE_PREF" fwmark "$table" table "$table"

  if [ "a$remote_dev" = "a" ] ; then
    ip route replace table "$table" default via "$remote_ip"
  else
    ip route replace table "$table" default via "$remote_ip" dev "$remote_dev"
  fi
}

cleanup_policy_route()
{
  while ip rule del pref "$DNSMASQ_ROUTE_RULE_PREF" 2>/dev/null
  do
    :
  done
}

setup_nftset_rules()
{
  local table="$1"
  local localnet_cidr="$2"
  local rule

  nft delete table ip "$DNSMASQ_ROUTE_NFT_TABLE" 2>/dev/null
  nft add table ip "$DNSMASQ_ROUTE_NFT_TABLE"
  nft "add set ip $DNSMASQ_ROUTE_NFT_TABLE $DNSMASQ_ROUTE_SET { type ipv4_addr; flags timeout; timeout $DNSMASQ_ROUTE_NFT_TIMEOUT_DEFAULT; }"
  nft "add chain ip $DNSMASQ_ROUTE_NFT_TABLE $DNSMASQ_ROUTE_NFT_CHAIN { type filter hook prerouting priority mangle; policy accept; }"

  rule="add rule ip $DNSMASQ_ROUTE_NFT_TABLE $DNSMASQ_ROUTE_NFT_CHAIN ip daddr @$DNSMASQ_ROUTE_SET"
  if [ "a$localnet_cidr" != "a" ] ; then
    rule="$rule ip saddr $localnet_cidr"
  fi
  rule="$rule meta mark set $table"
  nft "$rule"
}

cleanup_nftset_rules()
{
  nft delete table ip "$DNSMASQ_ROUTE_NFT_TABLE" 2>/dev/null
}

cleanup_runtime()
{
  local table="$1"

  rm -f "$DNSMASQ_ROUTE_CONF"
  cleanup_nftset_rules
  cleanup_policy_route "$table"
}

check_runtime_support()
{
  if ! command -v nft >/dev/null 2>&1 ; then
    echo "dnsmasq_route requires nftables and only supports OpenWrt 24+." >&2
    return 1
  fi
}

start()
{
  local remote_ip dns_server localnet_only table
  local sip netmask localnet_cidr

  remote_ip=$(uci -q get dnsmasq_route.config.remote_ip)
  dns_server=$(uci -q get dnsmasq_route.config.dns_server)
  localnet_only=$(uci -q get dnsmasq_route.config.localnet_only)
  table=$(uci -q get dnsmasq_route.config.table) #路由表编号

  if [ "a$localnet_only" = "a1" ] ; then
    sip=$(uci -q get network.lan.ipaddr)
    netmask=$(uci -q get network.lan.netmask)
    localnet_cidr="$sip/$(netmask_to_prefix "$netmask")"
  fi

  if [ "a$table" = "a" ] ; then
    table=107
  fi

  check_runtime_support || return 1

  if ! [ -e /etc/dnsmasq_route.ini ] ; then
    echo github.com > /etc/dnsmasq_route.ini
  fi

  cleanup_runtime "$table"

  setup_nftset_rules "$table" "$localnet_cidr"

  setup_policy_route "$table" "$remote_ip"
  echo "#auto gen from /etc/dnsmasq_route.ini" > "$DNSMASQ_ROUTE_CONF"
  grep -v -e '^#' -e '^$' /etc/dnsmasq_route.ini | \
  while read hostname note
  do
    echo "server=/$hostname/$dns_server" >> "$DNSMASQ_ROUTE_CONF"
    echo "nftset=/$hostname/4#ip#$DNSMASQ_ROUTE_NFT_TABLE#$DNSMASQ_ROUTE_SET" >> "$DNSMASQ_ROUTE_CONF"
  done

  /etc/init.d/dnsmasq restart
  echo "dnsmasq_route has started."
}

stop(){
  local table

  table=$(uci -q get dnsmasq_route.config.table)
  if [ "a$table" = "a" ] ; then
    table=107
  fi

  cleanup_runtime "$table"
  echo "dnsmasq_route has stopped."
  /etc/init.d/dnsmasq restart
}

restart(){
       start
}
